# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
childlogdir=/var/log/supervisor/
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:activemq]
directory=/usr/share/activemq
command=/usr/bin/activemq console
environment=ACTIVEMQ_HOME="/usr/share/activemq",ACTIVEMQ_CONFIG_DIR="/etc/activemq/instances-enabled/mco",ACTIVEMQ_USER=activemq,ACTIVEMQ_PIDFILE=/var/run/activemq/mco.pid
user=activemq
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
priority=600

[program:atd]
directory=/usr/sbin
command=/usr/sbin/atd
priority=100
startsecs=0
autostart=true
autorestart=false

[program:puppetdb]
directory=/opt/puppetlabs/
command=bash ./run_puppet_db.sh
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
priority=900

[program:mcollective]
directory=/opt/puppetlabs/puppet/bin/
command=/opt/puppetlabs/puppet/bin/ruby ./mcollectived --no-daemonize
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
priority=950

[program:puppet]
directory=/opt/puppetlabs
command=bash ./run_puppet.sh
priority=700
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

[program:redis]
directory=/usr/bin
command=/usr/bin/redis-server /etc/redis/redis.conf
priority=750
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

[program:postgresql]
command=/usr/lib/postgresql/9.4/bin/postgres -D /var/lib/postgresql/9.4/main -c config_file=/etc/postgresql/9.4/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
directory=/var/run/postgresql/

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler

